# ansible/templates/mnt-backups.automount
[Unit]
Description=backup automount

[Automount]
Where=/mnt/backups

[Install]
WantedBy=multi-user.target